# supervisor.mysql.conf
[program:mysql]
command=/usr/bin/mysqld_safe
