# supervisor.ssh.conf
[program:ssh]
command=/usr/sbin/sshd -D
