# supervisor.memcached.conf
[program:memcached]
command=/usr/bin/memcached -m 2000 -u nobody -l 0.0.0.0 -c 10240 -p 11211
